# deploy/kline_supervisor.conf
[program:web_kline]
command =python Main.py
user = root
directory=/webapps/web_kline
stdout_logfile=/webapps/web_kline/supervisor.log
autostart=true
autorestart=true